# files/dot_config/systemd/user/chroma.service
[Unit]
Description=Chroma server

[Service]
Type=simple
ExecStart=%h/.local/bin/chroma server

[Install]
WantedBy=default.target
